# bin/HBUpgradeCask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/joe/bin:/usr/local/opt/coreutils/libexec/gnubin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin:/opt/X11/bin:/Applications/Server.app/Contents/ServerRoot/usr/bin:/Applications/Server.app/Contents/ServerRoot/usr/sbin:/usr/local/MacGPG2/bin:/Applications/lb5pa2:/Users/joe/Library/Python/2.7/bin</string>
		<key>VAR</key>
		<string>VAL</string>
	</dict>
	<key>Label</key>
	<string>local.homebrew.upgrade-cask</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/joe/bin/brew-cask.sh</string>
		<string>upgrade</string>
	</array>
	<key>StandardErrorPath</key>
	<string>Homebrew-upgrade-cask.err</string>
	<key>StandardOutPath</key>
	<string>Homebrew-upgrade-cask.out</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>6</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/joe</string>
</dict>
</plist>
